# etc/profile.d/20-rvm.sh
if [[ -s "$HOME/.rvm/scripts/rvm" ]]; then
  # for .bash_profile
  source "$HOME/.rvm/scripts/rvm" # Load RVM into a shell session *as a function*
  # for .bashrc
  export PATH="$PATH:$HOME/.rvm/bin" # Add RVM to PATH for scripting
fi
