# etc/profile.d/README
10-*.sh  自前インストールや優先プログラムのPATH設定など
20-*.sh  10-*.shの設定に依存する設定
30-*.sh  その他
